# crontabs/root
0 4 * * 3  /scripts/gendns.sh && supervisorctl restart dnsmasq
