# files/crontabs/root
0 4 * * 3  cat /conf.d/main.conf | grep "pachost" | awk '{print $2}' | xargs sh /scripts/gendns.sh && supervisorctl restart dnsmasq
#*/1 * * * *  cat /conf.d/main.conf | grep "pachost" | awk '{print $2}' | xargs sh /scripts/gendns.sh && supervisorctl restart dnsmasq
