# supervisor.d/client.ini
[program:local]
command=/usr/bin/ss-local -c /conf.d/local.json
autostart=false

[program:tunnel]
command=/usr/bin/ss-tunnel -c /conf.d/tunnel.json
autostart=false

[program:redir]
command=/usr/bin/ss-redir -c /conf.d/redir.json
autostart=false

[program:dnsmasq]
command=/usr/sbin/dnsmasq -d
autostart=false
