# supervisor.d/client.ini
[program:local]
command=/usr/bin/ss-local -c /conf.d/local.json
autostart=false
