# supervisor.d/server.ini
[program:server]
command=/usr/bin/ss-server -c /conf.d/server.json
autostart=false
